# eyenotify.service
[Unit]
Description=Health Notification Service (Eye Protection)
After=graphical-session.target

[Service]
# Use simple type for long-running scripts
Type=simple

# %h expands to the user's home directory
ExecStart=/usr/bin/python3 %h/scripts/eyenotify.py

# Restart on failure with a short delay
Restart=on-failure
RestartSec=10

# Environment for GUI tools (notify-send) and unbuffered output
Environment=DISPLAY=:0
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
